# systemd/.config/systemd/user/redshift-on.service
[Unit]
Description=Redshift Service

[Service]
Type=simple
ExecStart=/usr/bin/redshift
ExecStop=/usr/bin/pkill redshift
RemainAfterExit=yes
