# systemd/.config/systemd/user/redshift-off.service
[Unit]
Description=Stop Redshift in the morning

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop redshift-on.service

